# frontend/cron/update_api_ip.crontab
* * * * * /home/ubuntu/update_api_ip.sh





